# backup/backup.service
[Unit]
Description=Monthly backup

[Service]
Type=oneshot
ExecStart=/home/anders/Documents/GitHub/config/backup/backup.sh
