# backup/backup.service
[Unit]
Description=Monthly backup

[Service]
Type=oneshot
ExecStart=/home/anders/Code/bash/backup.sh
User=anders
